# images/router/supervisord.conf
[supervisord]
logfile=/var/log/supervisor/supervisord.log 

[program:sshd]


[program:zebra]


[program:bgpd]


[program:ospfd]


[program:isisd]